# .config/autostart/Tracker.desktop
[Desktop Entry]
Type=Application
Name=Tracker
Description=Tracker
Exec=tracker daemon -k
